# etc/cron.d/magento.crontab

*/5 * * * * www-data /var/www/magento/cron.sh
